# avito-sync.service
[Unit]
Description=Avito Chats Auto Sync Service
After=network.target

[Service]
Type=simple
User=dannnnay
WorkingDirectory=/home/d/dannnnay/osagaming.store
Environment="SYNC_INTERVAL=60"
ExecStart=/usr/bin/python3 /home/d/dannnnay/osagaming.store/backend/auto_sync.py
Restart=always
RestartSec=10
StandardOutput=append:/home/d/dannnnay/osagaming.store/auto_sync.log
StandardError=append:/home/d/dannnnay/osagaming.store/auto_sync.log

[Install]
WantedBy=multi-user.target

